# youtube-downloader.service
# Systemd service file for YouTube Downloader
# Save this as /etc/systemd/system/youtube-downloader.service

[Unit]
Description=YouTube Downloader Web Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/youtube-downloader
Environment="PATH=/path/to/your/youtube-downloader/venv/bin"
ExecStart=/path/to/your/youtube-downloader/venv/bin/gunicorn --config gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
